# nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true

[program:nginx]
command=/usr/local/sbin/nginx -c /usr/local/nginx/conf/nginx.conf
autostart=true
