# mysql/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true

[program:mysqld]
command=/etc/init.d/mysqld start
autostart=true

